# profile.d/x-tools.sh
export PATH=$PATH:/opt/x-tools/arm-rpi-linux-gnueabihf/bin
